# site-cookbooks/appenv/templates/default/02_convert.conf.erb
[program:extract]
command=/var/www/application/current/app/Console/cake -app /var/www/application/current/app extract
user=<%= node['appenv']['OSS_BATCH_USER'] %>
autostart=true
autorestart=true
stdout_logfile=/var/log/extract.log
redirect_stderr=true
environment=OSS_RDS_HOSTNAME=<%= node['appenv']['OSS_RDS_HOSTNAME'] %>,OSS_RDS_USERNAME=<%= node['appenv']['OSS_RDS_USERNAME'] %>,OSS_RDS_PASSWORD=<%= node['appenv']['OSS_RDS_PASSWORD'] %>,OSS_RDS_DB_NAME=<%= node['appenv']['OSS_RDS_DB_NAME'] %>,OSS_BUCKET_NAME=<%= node['appenv']['OSS_BUCKET_NAME'] %>,OSS_IMAGE_BUCKET_NAME=<%= node['appenv']['OSS_IMAGE_BUCKET_NAME'] %>,OSS_USE_S3_STATIC_HOSTING=<%= node['appenv']['OSS_USE_S3_STATIC_HOSTING'] %>,OSS_REGION=<%= node['appenv']['OSS_REGION'] %>,OSS_SQS_URL=<%= node['appenv']['OSS_SQS_URL'] %><% if not node['appenv']['OSS_CDN_BASE_URL'].empty? %>,OSS_CDN_BASE_URL=<%= node['appenv']['OSS_CDN_BASE_URL'] %><% end %><% if not node['appenv']['OSS_AWS_ACCESS_ID'].empty? %>,OSS_AWS_ACCESS_ID=<%= node['appenv']['OSS_AWS_ACCESS_ID'] %><% end %><% if not node['appenv']['OSS_AWS_SECRET_KEY'].empty? %>,OSS_AWS_SECRET_KEY=<%= node['appenv']['OSS_AWS_SECRET_KEY'] %><% end %>
