# site-cookbooks/appenv/templates/default/02_convert.conf.erb
[program:extract]
command=/var/www/application/current/app/Console/cake -app /var/www/application/current/app extract
user=<%= node['appenv']['batch_user'] %>
autostart=true
autorestart=true
stdout_logfile=/var/log/extract.log
redirect_stderr=true
environment=RDS_HOSTNAME=<%= node['appenv']['RDS_HOSTNAME'] %>,RDS_USERNAME=<%= node['appenv']['RDS_USERNAME'] %>,RDS_PASSWORD=<%= node['appenv']['RDS_PASSWORD'] %>,RDS_DB_NAME=<%= node['appenv']['RDS_DB_NAME'] %>,BUCKET_NAME=<%= node['appenv']['BUCKET_NAME'] %>,IMAGE_BUCKET_NAME=<%= node['appenv']['IMAGE_BUCKET_NAME'] %>,USE_S3_STATIC_HOSTING=<%= node['appenv']['USE_S3_STATIC_HOSTING'] %>,REGION=<%= node['appenv']['REGION'] %>,SQS_URL=<%= node['appenv']['SQS_URL'] %>,CDN_BASE_URL=<%= node['appenv']['CDN_BASE_URL'] %><% if not node['appenv']['AWS_ACCESS_ID'].empty? %>,AWS_ACCESS_ID=<%= node['appenv']['AWS_ACCESS_ID'] %><% end %><% if not node['appenv']['AWS_SECRET_KEY'].empty? %>,AWS_SECRET_KEY=<%= node['appenv']['AWS_SECRET_KEY'] %><% end %>
